# deploy/nenya.service
# Nenya AI Gateway - systemd service unit
#
# Install:
#   curl -fsSL https://raw.githubusercontent.com/gumieri/nenya/main/install.sh | sudo sh
#   See: https://github.com/gumieri/nenya#quick-start
#
# Install paths:
#   Binary:    /usr/local/bin/nenya
#   Config:    /etc/nenya/  (directory, loaded and deep-merged)
#   Secrets:   /etc/nenya/secrets.json  (loaded via systemd credential)
#   Config ref: /etc/nenya/config.json.example (installed on first install if no config exists)
#
# Hot reload:
#   systemctl reload nenya

[Unit]
Description=Nenya AI Gateway & Bouncer
Requires=nenya.socket
After=nenya.socket

[Service]
Type=simple
DynamicUser=yes 
ExecStart=/usr/local/bin/nenya
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

LoadCredential=secrets:/etc/nenya/secrets.json

[Install]
WantedBy=multi-user.target
