# deploy/nenya.service
# Nenya AI Gateway - systemd service unit
#
# Install paths (adjust if deploying elsewhere):
#   Binary:    /usr/local/bin/nenya
#   Config:    /etc/nenya/config.json
#   Secrets:   /etc/nenya/secrets.json  (loaded via systemd credential)
#   Install:   sudo mise run install
#   Enable:    sudo systemctl enable --now nenya

[Unit]
Description=Nenya AI Gateway & Bouncer
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes 
ExecStart=/usr/local/bin/nenya
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

LoadCredential=secrets:/etc/nenya/secrets.json

[Install]
WantedBy=multi-user.target
